# s6/services.d/flux-waf/finish
#!/bin/sh
echo "[s6-flux-waf] Service exited. Exit code: $1, signal: $2"
